# deploy/schrubles.socket
[Unit]
Description=paracetaotimo socket

[Socket]
ListenStream=/run/paracetaotimo/socket

[Install]
WantedBy=sockets.target
